# deploy/systemd/ulanzi-web.service
[Unit]
Description=Ulanzi static webapp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ulanzi-bridge
ExecStart=/usr/bin/python3 -m bridge.webapp_server --host 0.0.0.0 --port 80 --root /var/www/ulanzi
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
